# autoimport-start-listener.conf
author "Nathan West <nathan@mimedia.com>"
description "Start the file-bridge service on the incoming directory"

start on filesystem
task

exec start file-bridge FILE=/srv/apt/incoming/